# debian/fancontrol/lib/systemd/system/fancontrol.service
[Unit]
Description=Custom Fan Control Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/fancontrol
Restart=always

[Install]
WantedBy=multi-user.target